# systemd/iot.service
[Unit]
Description=IOT Service

[Service]
Type=simple
ExecStart=/opt/iotproject/bin/webservice.bash

[Install]
WantedBy=multi-user.target